# deploy/templates/systemd/yodeploy-flask.service.template
[Unit]
Description=yodeploy flask server

[Service]
ExecStart=/usr/local/bin/yodeploy-server
Restart=on-failure

[Install]
WantedBy=multi-user.target
